# projects/sdr_receiver/multicast/sdr-receiver.conf
description "Red Pitaya SDR server"

start on runlevel [345]
stop on runlevel [^345]

respawn

exec /usr/bin/sdr-receiver
