# projects/sdr_transceiver_ft8/app/ft8.cron
* * * * * cd /dev/shm && /media/mmcblk0p1/apps/sdr_transceiver_ft8/decode-ft8.sh >> decode-ft8.log 2>&1 &
*/6 * * * * cd /dev/shm && /media/mmcblk0p1/apps/sdr_transceiver_ft8/upload-ft8.sh >> upload-ft8.log 2>&1 &
#* * * * * cd /dev/shm && /media/mmcblk0p1/apps/common_tools/update-corr.sh 125 >> update-corr.log 2>&1 &
